# com.skybg.cache-cleanup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.skybg.cache-cleanup</string>
  <key>ProgramArguments</key>
  <array>
    <string>__PROJECT_DIR__/bin/skybg-cleanup</string>
  </array>
  <key>StartInterval</key>
  <integer>300</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardErrorPath</key>
  <string>__LOG_DIR__/cache-cleanup.log</string>
</dict>
</plist>
